# docker/rootfs/etc/s6-overlay/s6-rc.d/swift-main/run
#!/command/with-contenv /bin/sh

exec 2>&1
exec s6-setuidgid swift swift-init main start --no-daemon
