# docker/rootfs/etc/services.d/swift-main/run
#!/usr/bin/execlineb -P

with-contenv

exec s6-setuidgid swift /usr/local/bin/swift-init main start --no-daemon
