# docker/rootfs/etc/s6-overlay/s6-rc.d/memcached/run
#!/bin/sh

exec 2>&1
exec memcached -u root -l 127.0.0.1
