# docker/rootfs/etc/s6-overlay/s6-rc.d/keystone/run
#!/command/with-contenv /bin/sh

exec 2>&1
exec s6-setuidgid keystone keystone-wsgi-public -b 0.0.0.0 -p 5000
